# full-sys/etc/crontabs/root
0 * * * * /root/backup.sh >> /root/openwrt-backups.log 2>&1
